# config/autostart/auto-brightness.desktop
[Desktop Entry]
Name=auto-brightness
Exec=bash -c "source ~/.dotfiles/bin/auto-brightness.sh"
Type=Application
