# config/autostart/auto-brightness.desktop
[Desktop Entry]
Name=auto-brightness
Exec=bash -c "~/.dotfiles/bin/auto-brightness &> /dev/null"
Type=Application
