# config/autostart/set-resolution-remote.desktop
[Desktop Entry]
Name=set-resolution-remote
Exec=bash -c "~/.dotfiles/bin/set-resolution-remote &> /dev/null"
Type=Application
